# shm_service/shm-service.service
[Unit]
Description=RaspiPLC Shared Memory Service
Documentation=https://github.com/ASolchen/RaspiPLC
After=raspiplc-shm-init.service
Requires=raspiplc-shm-init.service

[Service]
Type=simple

# Run as root for now (needed for /dev/shm and /run socket ownership)
User=root
Group=root

# Environment
Environment=RASPIPLC_TAGS_FILE=/etc/raspiplc/tags.json

# Executable
ExecStart=/usr/bin/env python3 /opt/raspiplc/shm-service/shm_service.py

# Restart behavior
Restart=always
RestartSec=1

# Hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

